# cookbooks/wendy/files/default/etc/cron.d/mig-loader
30 */2 * * * root /sbin/mig-loader >/dev/null 2>&1
